# runtime/etc/profile.d/codexpad.sh
export CODEX_HOME="${CODEX_HOME:-/root/.codex}"

cd /root/workspace 2>/dev/null || true
